# server/utils/storage/keep-storage-mounted.timer
[Unit]
Description=Run Keepalive HDD script periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
